# crontab
0 2 * * * root /auto_dump.sh >> /auto_dump_log.txt 2>&1
